# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stdout_logfile=/var/log/supervisor/php-fpm.out.log
priority=999

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=998

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/cron.err.log
stdout_logfile=/var/log/supervisor/cron.out.log
priority=997

[program:heartbeat]
command=/bin/sh -lc 'if [ "${HEARTBEAT_ENABLED:-1}" = "1" ]; then exec /usr/local/bin/php /var/www/html/heartbeat.php; else exec tail -f /dev/null; fi'
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/var/log/supervisor/heartbeat.err.log
stdout_logfile=/var/log/supervisor/heartbeat.out.log
priority=996

[program:crm-user-receiver]
command=/usr/local/bin/php /var/www/html/crm_user_receiver.php
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/var/log/supervisor/crm-user-receiver.err.log
stdout_logfile=/var/log/supervisor/crm-user-receiver.out.log
priority=995

[program:crm-company-receiver]
command=/usr/local/bin/php /var/www/html/crm_company_receiver.php
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/var/log/supervisor/crm-company-receiver.err.log
stdout_logfile=/var/log/supervisor/crm-company-receiver.out.log
priority=994

[program:kassa-batch-receiver]
command=/usr/local/bin/php /var/www/html/kassa_batch_receiver.php
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/var/log/supervisor/kassa-batch-receiver.err.log
stdout_logfile=/var/log/supervisor/kassa-batch-receiver.out.log
priority=993

[program:kassa-invoice-receiver]
command=/usr/local/bin/php /var/www/html/kassa_invoice_requested_receiver.php
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/var/log/supervisor/kassa-invoice-receiver.err.log
stdout_logfile=/var/log/supervisor/kassa-invoice-receiver.out.log
priority=992
