# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=999

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=998

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=997

[program:heartbeat]
command=/bin/sh -lc 'if [ "${HEARTBEAT_ENABLED:-1}" = "1" ]; then exec /usr/local/bin/php /var/www/html/heartbeat.php; else exec tail -f /dev/null; fi'
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=996

[program:crm-user-receiver]
command=/usr/local/bin/php /var/www/html/crm_user_receiver.php
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=995

[program:crm-company-receiver]
command=/usr/local/bin/php /var/www/html/crm_company_receiver.php
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=994

[program:kassa-batch-receiver]
command=/usr/local/bin/php /var/www/html/kassa_batch_receiver.php
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=993

[program:kassa-invoice-receiver]
command=/usr/local/bin/php /var/www/html/kassa_invoice_requested_receiver.php
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=992

[program:status-check]
command=/usr/local/bin/php /var/www/html/status_check.php
autostart=true
autorestart=true
startsecs=0
startretries=999
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=991

